# supervisord.conf
[program:logentries]
command = /usr/bin/python /home/dotcloud/data/le monitor
